# roles/linux/templates/unit.service.j2
[Unit]
Description={{ description }}

[Service]
Type={{ type }}
ExecStart={{ exec_start }}

[Install]
WantedBy={{ wanted_by }}
